# image/etc/cron.d/sync
* * * * * root  flock -x -n -E 0 /var/run/sync.lock /opt/build/bin/sync.sh | logger -t syncer -i -p local5.warn
